# deploy/systemd/marketlab-youtube.service
[Unit]
Description=Market Sentiment Lab - YouTube pipeline
After=network-online.target

[Service]
Type=simple
User=gaston
WorkingDirectory=/home/gaston/projects/market-sentiment-lab
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/gaston/projects/market-sentiment-lab/lab/src
ExecStart=/usr/bin/timeout 30m /bin/bash -lc 'source lab/.venv/bin/activate; export OPENAI_API_KEY="$(sed -n "s/^OPENAI_API_KEY=//p" /home/gaston/apps/openclaw/.env | head -n1)"; ./lab/scripts/run_youtube_pipeline.sh'
Restart=on-failure
RestartSec=20
StandardOutput=append:/home/gaston/logs/marketlab/systemd.log
StandardError=append:/home/gaston/logs/marketlab/systemd.err

[Install]
WantedBy=multi-user.target
